# secondary-droplet/systemd/youtube-fallback.service.d/override.conf
[Service]
EnvironmentFile=
EnvironmentFile=-/etc/youtube-fallback.env
