# secondary-droplet/systemd/yt-restapi.service
[Unit]
Description=BWB Status Monitor (porta 8080)
After=network.target

[Service]
User=yt-restapi
Group=yt-restapi
EnvironmentFile=/etc/yt-restapi.env
PermissionsStartOnly=true
ExecStartPre=/bin/sh -c '
    touch /run/youtube-fallback.mode && \
    chown yt-restapi:yt-restapi /run/youtube-fallback.mode && \
    chmod 0644 /run/youtube-fallback.mode
'
ExecStart=/usr/bin/python3 /usr/local/bin/bwb_status_monitor.py
Restart=always
RestartSec=3
# Nota: o monitor controla youtube-fallback via sudo; este serviço necessita
# de permitir a escalada pontual, caso contrário sudo falha com "no new
# privileges" e o fallback nunca arranca.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
RuntimeDirectory=youtube-fallback
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
