# secondary-droplet/systemd/ytc-web-backend.service
[Unit]
Description=YTC Web Backend (YouTube status API)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ytc-web-service/app
EnvironmentFile=/etc/ytc-web-backend.env
Environment=PYTHONUNBUFFERED=1
Environment=YTC_WEB_BACKEND_HOST=127.0.0.1
Environment=YTC_WEB_BACKEND_PORT=8081
ExecStart=/opt/ytc-web-service/venv/bin/uvicorn app:app --host ${YTC_WEB_BACKEND_HOST} --port ${YTC_WEB_BACKEND_PORT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
